# ziman-cleaner.service
[Unit]
Description=ZiVPN Password Cleaner

[Service]
Type=oneshot
ExecStart=/etc/zivpn/ziman-cleaner.sh
